# onboot.service
[Unit]
Description=My script

[Service]
ExecStart=/root/onboot.sh

[Install]
WantedBy=multi-user.target 
